# containers/nginx/supervisor/socketio_server.conf
[program:socketio]
command = /usr/bin/node /var/www/laravel/socket.js
directory = /var/www/laravel
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
